# megg-iot-backend.service
[Unit]
Description=MEGG IoT Backend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=megg
Group=megg
WorkingDirectory=/home/megg/megg-iot-backend
Environment=ARDUINO_PORT=/dev/ttyUSB0
Environment=ARDUINO_BAUDRATE=9600
Environment=IOT_BACKEND_HOST=0.0.0.0
Environment=IOT_BACKEND_PORT=8765
ExecStart=/home/megg/megg-iot/venv/bin/python simple_iot_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give the service access to serial ports
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
